# config/supervisord/app.ini
[program:app_supervisord]
directory=/app
command=node index.js
stderr_logfile=/var/log/app_supervisord.err.log
stdout_logfile=/var/log/app_supervisord.out.log